# system/etc/systemd/system/methaneos-update.timer
[Unit]
Description=Run MethaneOS Update Check Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
