# system/etc/systemd/system/methaneos-update.service
[Unit]
Description=Check for MethaneOS Updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/methaneos-update check
